# data/release/systemd/goose-proxy.service
[Unit]
Description=Goose Proxy Daemon Service
Documentation=man:goose-proxy(7)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
RemainAfterExit=no
ExecStart=/usr/bin/goose-proxy

# Run as an ephemeral unprivileged user
DynamicUser=true

# Filesystem: read-only except /dev, /proc, /sys
ProtectHome=true
ProtectSystem=strict

# Load RHSM certificates via systemd credentials so the DynamicUser can
# access files that are owned by root in /etc/pki/consumer.
LoadCredential=cert.pem:/etc/pki/consumer/cert.pem
LoadCredential=key.pem:/etc/pki/consumer/key.pem

# Private /tmp so the service can't snoop on other processes
PrivateTmp=true

# Device access: only /dev/null, /dev/urandom, /dev/random
PrivateDevices=true
DeviceAllow=/dev/null rw
DeviceAllow=/dev/urandom r
DeviceAllow=/dev/random r

# Kernel and system protection
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectClock=true
ProtectProc=invisible
ProcSubset=pid
ProtectHostname=true

# Network: restrict to IPv4, IPv6, and Unix sockets
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Namespace and privilege restrictions
RestrictNamespaces=true
CapabilityBoundingSet=
NoNewPrivileges=true
PrivateUsers=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true

# System call filtering: only native arch, only basic service calls
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Memory protection
MemoryDenyWriteExecute=true

# Resource limits
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target